# templates/gpb-microservice.service.erb
[Unit]
Description=<%= @name %> microservice
After=nats.service redis.service
Requires=docker.service

[Service]
User=<%= @user %>
Restart=always
Environment=NATS_URI=<%= @natsuri %>
WorkingDirectory=/opt/go/src/github.com/<%= @organization %>/<%= @name %>
ExecStartPre=/usr/bin/timeout 120 /bin/bash -c 'while ! nc -z -v -q 1 localhost 4222; do sleep 1; done'
ExecStartPre=/usr/bin/timeout 120 /bin/bash -c 'while ! nc -z -v -q 1 localhost 6379; do sleep 1; done'
ExecStart=/opt/go/bin/<%= @name %>

[Install]
WantedBy=multi-user.target
