# templates/redis.service.erb
[Unit]
Description=redis
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStartPre=-/usr/bin/docker stop redis
ExecStartPre=-/usr/bin/docker rm redis
ExecStartPre=/usr/bin/docker pull redis
ExecStart=/usr/bin/docker run --name redis -p 6379:6379 -v /etc/redis/redis.conf:/usr/local/etc/redis/redis.conf -v <%= @datavol %>:/data redis redis-server /usr/local/etc/redis/redis.conf

[Install]
WantedBy=multi-user.target
